# torchpod_root/etc/services.d/sshd/run
#!/command/with-contenv bash

set -eu
set -o pipefail

echo "GEMFIELD: start opensshd..."
exec /usr/sbin/sshd -D
