# torchpod_root/etc/services.d/dbus/run
#!/command/with-contenv bash

set -eu
set -o pipefail

exec /usr/bin/dbus-daemon --system --nofork
